# board/opendingux/gcw0/target_skeleton/etc/udev/rules.d/61-network.rules
SUBSYSTEM!="net", GOTO="net_end"

ACTION=="add", RUN+="/sbin/ifup %k"
ACTION=="remove", RUN+="/sbin/ifdown %k"

LABEL="net_end"
