# board/opendingux/od-rescue-cd/target_skeleton/etc/inittab
null::sysinit:/bin/mount -t proc proc /proc
tty1::sysinit:/bin/sh /root/opendingux-bootdisc.sh
null::shutdown:/bin/umount -a -r
