# board/opendingux/od-rescue-cd/target_skeleton/etc/inittab
null::sysinit:/bin/mount -t proc proc /proc
null::sysinit:/bin/mount -t sysfs sysfs /sys
tty1::sysinit:/bin/sh /usr/share/od-rescue-cd/opendingux-bootdisc.sh
null::shutdown:/bin/umount -a -r
